# systemd/unitysphere-improved.service
[Unit]
Description=UnitySphere Django Application
Documentation=https://github.com/your-org/unitysphere
After=docker.service
Requires=docker.service
PartOf=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/unitysphere
User=unitysphere
Group=unitysphere

# Environment
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/opt/unitysphere/.env

# Pre-start checks
ExecStartPre=/usr/bin/docker compose -f /opt/unitysphere/docker-compose.yaml config -q
ExecStartPre=-/usr/bin/docker compose -f /opt/unitysphere/docker-compose.yaml down

# Start application
ExecStart=/usr/bin/docker compose -f /opt/unitysphere/docker-compose.yaml up -d

# Health check after start
ExecStartPost=/bin/sleep 10
ExecStartPost=/bin/bash -c 'curl -f http://localhost:8001/api/v1/ai/health/ || exit 1'

# Reload
ExecReload=/usr/bin/docker compose -f /opt/unitysphere/docker-compose.yaml restart fnclub

# Stop
ExecStop=/usr/bin/docker compose -f /opt/unitysphere/docker-compose.yaml down

# Restart policy
Restart=on-failure
RestartSec=10s
StartLimitInterval=200
StartLimitBurst=3

# Security
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=unitysphere

[Install]
WantedBy=multi-user.target
